# mailhog.conf
description "MailHog"

start on (net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 10 5

script
	exec /root/MailHog -smtp-bind-addr 0.0.0.0:25
end script

